# deploy/etc/systemd/system/elasticsearch.service.d/limit.conf
[Service]
LimitMEMLOCK=unlimited
LimitNOFILE=65536
LimitNPROC=4096
LimitAS=infinity
LimitFSIZE=infinity
